# mdm/com.aerowholesale.qc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aerowholesale.qc</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            # Prevent sleep
            caffeinate -dis &amp;
            CAFFEINE_PID=$!

            # Start Flask QC server
            cd /Users/aeroserver/aw-qc-system
            /usr/bin/python3 run.py

            kill $CAFFEINE_PID 2>/dev/null
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>UserName</key>
    <string>aeroserver</string>
    <key>WorkingDirectory</key>
    <string>/Users/aeroserver/aw-qc-system</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FLASK_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/aw-qc-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aw-qc-server.err</string>
</dict>
</plist>
